# sdk_container/src/third_party/coreos-overlay/app-emulation/open-vm-tools/files/vmtoolsd.service
[Unit]
Description=Service for virtual machines hosted on VMware
Documentation=https://github.com/vmware/open-vm-tools
ConditionVirtualization=vmware

[Service]
# Flatcar: Terrible.
ExecStartPre=/bin/bash -c "ln -srf /usr/share/vmtoolsd/pam.d/vmtoolsd /etc/pam.d/vmtoolsd"
ExecStart=/usr/bin/vmtoolsd
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
